# devel/etc/catkin/profile.d/10.rosbuild.sh
# generated from rosbuild/env-hooks/10.rosbuild.sh.em

# env variables in develspace
export ROS_ROOT="/Users/adityabhat/Downloads/src/ros/rosbuild"
