# devel/etc/catkin/profile.d/15.rosbash.bash
# generated from rosbash/env-hooks/15.rosbash.bash.em

. "/Users/adityabhat/Downloads/src/ros/rosbash/rosbash"
