# devel/etc/catkin/profile.d/15.rosbash.fish
# generated from rosbash/env-hooks/15.rosbash.fish.em

. "/Users/adityabhat/Downloads/src/ros/rosbash/rosfish"
