# devel/etc/catkin/profile.d/15.rosbash.zsh
# generated from rosbash/env-hooks/15.rosbash.zsh.em

. "/Users/adityabhat/Downloads/src/ros/rosbash/roszsh"
