# devel/etc/catkin/profile.d/15.rosbash.tcsh
# generated from rosbash/env-hooks/15.rosbash.tcsh.em

. "/Users/adityabhat/Downloads/src/ros/rosbash/rostcsh"
